# gpu-keepalive-simple.service
[Unit]
Description=Simple GPU Keep-Alive Service (No pynvml)
After=multi-user.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/shci-app
ExecStart=/usr/bin/python3 /var/www/shci-app/gpu-keepalive-simple.py --interval 30 --gpu-id 0
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
